# powershell/profile.d/00-env.ps1
# ---- 環境変数・共通フラグ ----

# 例: 共通環境変数
# $env:EDITOR = "code"
# $env:PAGER = "less"

# 例: フラグ設定
# $global:EnableVerboseLogging = $false
