# powershell/profile.d/10-github-copilot.ps1
function Start-CopilotYolo {
    copilot --allow-all-tools @args
}

Set-Alias -Name copilot-yolo -Value Start-CopilotYolo

