# powershell/profile.d/20-claude-code.ps1
function Start-ClaudeYolo {
    & claude --dangerously-skip-permissions @args
}

Set-Alias -Name claude-yolo -Value Start-ClaudeYolo